# .config/systemd/user/mako.service
[Unit]
Description=Lightweight Notification Daemon for Wayland

[Service]
ExecStart=/usr/bin/mako

[Install]
WantedBy=default.target
